# komodo/ansible/roles/komodo_periphery/templates/periphery.service.j2
# =============================================================================
# KOMODO PERIPHERY SYSTEMD SERVICE
# =============================================================================
# Generated by Ansible - Do not edit manually
# Template: periphery.service.j2
# =============================================================================

[Unit]
Description=Komodo Periphery Agent
Documentation=https://komo.do/docs
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User={{ periphery_user }}
Group={{ periphery_group }}
ExecStart={{ periphery_binary_path }} --config-path {{ periphery_config_dir }}/periphery.config.toml
Restart=on-failure
RestartSec=10
TimeoutStartSec=30

# -----------------------------------------------------------------------------
# ENVIRONMENT
# -----------------------------------------------------------------------------
Environment="RUST_LOG={{ periphery_log_level }}"
{% if periphery_logging_pretty %}
Environment="PERIPHERY_LOGGING_PRETTY=true"
{% endif %}
{% for secret in periphery_agent_secrets | default([]) %}
Environment="{{ secret.name }}={{ secret.value }}"
{% endfor %}

# -----------------------------------------------------------------------------
# SECURITY HARDENING
# -----------------------------------------------------------------------------
NoNewPrivileges=no
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths={{ periphery_data_dir }}
ReadWritePaths=/var/run/docker.sock

# -----------------------------------------------------------------------------
# RESOURCE LIMITS
# -----------------------------------------------------------------------------
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
